# docs/m7/deployment/systemd/soa-runner-crl-refresh.service
# CRL refresh helper for the Runner (reference).
#
# Core §7.3.1 CRL cache has a 2-hour stale ceiling — fetch MUST happen before
# the ceiling or /ready flips to 503 {crl-stale}. The Runner's in-process
# scheduler (BootOrchestrator.refreshIntervalMs) handles this in nominal
# conditions; this timer exists as a belt-and-suspenders backup that pokes
# the Runner's POST /crl/refresh operator endpoint on a fixed cadence.
#
# Requires the Runner to be wired with opts.crlRefresh set (the reference
# start-runner wires this automatically when SOA_RUNNER_BOOTSTRAP_BEARER is
# set). The endpoint requires admin-bearer auth — SOA_OPERATOR_BEARER must
# be provided via the EnvironmentFile below, which operators typically
# populate from a secrets manager at provisioning time.
#
# NOTE: Do NOT enable this service directly. It is fired by its matching
# .timer unit.

[Unit]
Description=SOA-Harness Runner — CRL refresh trigger (§7.3.1 backup)
Requires=soa-runner.service
After=soa-runner.service

[Service]
Type=oneshot
User=soa
Group=soa
EnvironmentFile=/etc/soa/runner.env
ExecStart=/usr/bin/curl -fsS --max-time 30 \
  -H "Authorization: Bearer ${SOA_OPERATOR_BEARER}" \
  -X POST http://127.0.0.1:7700/crl/refresh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=soa-runner-crl-refresh
